# pgsql/pgpro/systemd/backup-pgsql.timer
[Unit]
Description=Backup and maintenance PosrgeSQL DB

[Timer]
OnCalendar=*-*-* 01:15:00

[Install]
WantedBy=timers.target
